# crontab
* * * * * cd /myapp/lib/tasks && bundle exec rake notify_user_mailer:notify_user_mail